# src/autostart/disableInputs.desktop
[Desktop Entry]
Name=noInputs
Exec=lxterminal -e "/usr/bin/python3 /home/pi/Desktop/proyectoFInalEmbebidos/src/disableInputs.py"
Type=Application

