# src/autostart/mednaffe.desktop
[Desktop Entry]
Name=mednaffe
Exec=lxterminal -e "mednaffe"
Type=Application
