# src/autostart/qjoypad.desktop

[Desktop Entry]
Name=qjoypad
Exec=lxterminal -e "qjoypad 'Rasp'"
Type=Application
